# misc/init.d.example
#! /bin/sh

### BEGIN INIT INFO
# Provides:          jamboree
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the jamboree server
# Description:       starts jamboree using start-stop-daemon
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=jb
NAME=jamboree
DESC=jamboree

#test -x $DAEMON || exit 0

# Include jamboree defaults if available
if [ -f /etc/default/jamboree ] ; then
  . /etc/default/jamboree
fi

set -e

case "$1" in
  start)
    echo -n "Starting $DESC: "
    $DAEMON start
    ;;
  stop)
    echo -n "Stopping $DESC: "
    $DAEMON stop
    ;;
  restart|force-reload)
    echo -n "Restarting $DESC: "
    $DAEMON stop
    $DAEMON start
    ;;
  *)
    N=/etc/init.d/$NAME
    echo "Usage: $N {start|stop|restart}" >&2
    exit 1
    ;;
esac

exit 0
